# php7-fpm/crontab
*/5 * * * * root /usr/local/bin/php /code/console.php update >> /var/log/cron.log 2>&1
